# systemd/systemd/user/clipmenu.service
[Unit]
Description=Starts clipmenu daemon

[Service]
ExecStart=/usr/bin/clipmenud
Restart=always

[Install]
WantedBy=default.target
